# claude-wrapper-service/claude-wrapper.service
[Unit]
Description=Claude Analysis Wrapper Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tbticket/claude-wrapper-service
ExecStart=/usr/local/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
